# extras/systemd/glustereventsd.service.in
[Unit]
Description=Gluster Events Notifier
After=syslog.target network.target

[Service]
Environment=PYTHONPATH=@BUILD_PYTHON_SITE_PACKAGES_EXPANDED@:$PYTHONPATH
Type=simple
ExecStart=@SBIN_DIR@/glustereventsd
ExecReload=/bin/kill -SIGUSR2 $MAINPID
KillMode=control-group

[Install]
WantedBy=multi-user.target
