# extras/systemd/glusterd-mgmt.service.in
[Unit]
Description=GlusterFS Management Daemon
After=network.target rpcbind.service
Before=network-online.target

[Service]
Type=forking
PIDFile=@localstatedir@/lib/glusterd/run/glusterd.pid
LimitNOFILE=65536
ExecStart=@prefix@/sbin/glusterd -p @localstatedir@/lib/glusterd/run/glusterd.pid
KillMode=process
Restart=always

[Install]
WantedBy=multi-user.target
