# extras/systemd/glusterfssharedstorage.service.in
[Unit]
Description=Mount glusterfs sharedstorage
Requires=glusterd.service remote-fs-pre.target local-fs.target

[Service]
Type=forking
ExecStart=@GLUSTERFS_LIBEXECDIR@/mount-shared-storage.sh
Restart=on-failure
RestartSec=3
RestartForceExitStatus=1

[Install]
WantedBy=multi-user.target
