# extras/systemd/glusterd.service.in
[Unit]
Description=GlusterFS, a clustered file-system server
After=network.target rpcbind.service
Before=network-online.target

[Service]
Type=forking
PIDFile=/run/glusterd.pid
LimitNOFILE=65536
ExecStart=@prefix@/sbin/glusterd -p /run/glusterd.pid
KillMode=process

[Install]
WantedBy=multi-user.target
